# ffl-testing@.socket
[Unit]
Description=ffl-testing on port %i

[Socket]
ListenStream=127.0.0.1:%i

[Install]
WantedBy=sockets.target
